# service/dirpicsubscriber.service
[Unit]
Description=DirPic-Subscriber
After=network.target

[Service]
Type=notify
WorkingDirectory=/home/dirpic/runtime
RootDirectory=/home/dirpic/root
User=dirpic
Group=dirpic
ExecStart=/usr/bin/dirpicsubscriber -root=/home/dirpic/runtime
ExecReload=/bin/kill -USR1 $MAINPID
BindReadOnlyPaths=/usr/bin:/usr/bin
BindReadOnlyPaths=/usr/lib:/usr/lib
Environment=HOME=/home/dirpic/env

[Install]
WantedBy=multi-user.target
